# conf/supervisord.d/api_server.ini
[program:api_server]
command = /srv/urlist/python-env/run.sh /srv/urlist/motherbrain/api_server/server.py --port=9090
directory = /srv/urlist/motherbrain/api_server/
stderr_logfile = /var/log/urlist/api_server-stderr.log
stdout_logfile = /var/log/urlist/api_server-stdout.log
numprocs = 1
autostart=true
