# conf/supervisord.d/cloudsync.ini
[program:cloudsync]
command = /srv/urlist/go/bin/run.sh cloudsync
enviroment = GOPATH=/srv/urlist/go,PATH=/srv/urlist/bin
directory = /srv/urlist/go
stderr_logfile = /var/log/urlist/cloudsync-stderr.log
stdout_logfile = /var/log/urlist/cloudsync-stdout.log
numprocs = 1
autostart=true
