# conf/supervisord.d/embed.ini
[program:embed]
directory=/srv/urlist/go
enviroment = GOPATH=/srv/urlist/go
command = /srv/urlist/go/bin/run.sh embed
stderr_logfile = /var/log/urlist/embed-stderr.log
stdout_logfile = /var/log/urlist/embed-stdout.log
numprocs = 1
autostart=true
