# conf/supervisord.d/favicon.ini
[program:favicon]
command = /srv/urlist/go/bin/run.sh favicon
enviroment = GOPATH=/srv/urlist/go
directory = /srv/urlist/go
stderr_logfile = /var/log/urlist/favicon-stderr.log
stdout_logfile = /var/log/urlist/favicon-stdout.log
numprocs = 1
autostart=true
