# conf/supervisord.d/fetch.ini
[program:fetch]
command = /srv/urlist/go/bin/run.sh fetch
enviroment = GOPATH=/srv/urlist/go
directory = /srv/urlist/go
stderr_logfile = /var/log/urlist/fetch-stderr.log
stdout_logfile = /var/log/urlist/fetch-stdout.log
numprocs = 1
autostart=true
