# conf/supervisord.d/seo.ini
[program:seo]
directory=/srv/urlist/go/
enviroment = GOPATH=/srv/urlist/go/
command = /srv/urlist/go/bin/run.sh seo
stderr_logfile = /var/log/urlist/seo-stderr.log
stdout_logfile = /var/log/urlist/seo-stdout.log
numprocs = 1
autostart=true
