# conf/supervisord.d/web_client.ini
[program:web_client]
command = /srv/urlist/python-env/run.sh /srv/urlist/webclient/server/server.py
enviroment = PYTHONPATH="/srv/urlist/webclient/"
directory = /srv/urlist/webclient/server
stderr_logfile = /var/log/urlist/webclient-stderr.log
stdout_logfile = /var/log/urlist/webclient-stdout.log
numprocs = 1
autostart=true
