# conf/supervisord.d/widget.ini
[program:widget]
command = /srv/urlist/go/bin/run.sh widget
enviroment = GOPATH=/srv/urlist/go
directory = /srv/urlist/go
stderr_logfile = /var/log/urlist/widget-stderr.log
stdout_logfile = /var/log/urlist/widget-stdout.log
numprocs = 1
autostart=true
