# conf/supervisord.d/workers.ini
[program:workers]
command = /srv/urlist/python-env/run.sh /srv/urlist/motherbrain/workers
enviroment = PYTHONPATH=/srv/urlist/
directory = /srv/urlist/motherbrain/
stderr_logfile = /var/log/urlist/workers-stderr.log
stdout_logfile = /var/log/urlist/workers-stdout.log
numprocs = 1
autostart=true
killasgroup=true
